# etc.init.ecospec.conf
# Info
description "EcoSpec"
author      "Donald J. Bales"

# Events
start on startup
stop on shutdown

# Automatically respawn
respawn
respawn limit 20 5

# Run the script!
# Note, in this example, if your script returns
# the string "ERROR", the daemon will stop itself.
script
    [ $(exec /usr/bin/python /home/ecospec/ecospec.py) = 'ERROR' ] && ( stop; exit 1; )
end script
